# examples/sample-evidence-realistic/sigma-rules/suspicious_scheduled_task.yml
title: Scheduled Task Created Shortly After Unusual USB Insertion
id: 4f3e21aa-1234-4c9d-9b2e-dart-002
description: Catches scheduled tasks authored by SYSTEM that reference user
  Download directories — a common remote-hands persistence pattern.
status: experimental
author: Agentic-DART
level: medium
logsource:
  product: windows
  service: task_scheduler
detection:
  selection:
    event_type: task_registration
    author: NT AUTHORITY\SYSTEM
    command|contains:
      - \Users\
      - \Downloads\
  condition: selection
tags:
  - attack.persistence
  - attack.t1053
